# systemd-units/media-ssd120.automount
[Unit]
Description=Automount Ssd nvme 1 tera

[Automount]
Where=/media/ssd120
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
