# systemd-units/media-games.automount
[Unit]
Description=Automount Ssd nvme 1 tera

[Automount]
Where=/media/games
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
